# layout/DEBIAN/postinst
#!/bin/sh

# Post-installation script for ProjectX

# Mark script execution
echo "Running WeaponX post-installation script..."
echo "===================================="

# Detect rootless or traditional jailbreak
if [ -f /Library/LaunchDaemons/wangmc.projectx.plist ]; then
    JBPREFIX=""
    echo "✅ Traditional jailbreak detected, using no prefix"
    if [ -f /Library/LaunchDaemons/wangmc.projectx.plist ]; then
        launchctl load /Library/LaunchDaemons/wangmc.projectx.plist || true
    fi
    mv "/Library/ProjectXTweak.plist" "/Library/MobileSubstrate/DynamicLibraries/ProjectXTweak.plist"
elif [ -d /var/jb ]; then
    JBPREFIX="/var/jb"
    echo "✅ Rootless jailbreak detected, using prefix: $JBPREFIX"
    if [ -f /var/jb/Library/LaunchDaemons/wangmc.projectx.plist ]; then
        launchctl load /var/jb/Library/LaunchDaemons/wangmc.projectx.plist || true
    fi
    mv "/var/jb/Library/ProjectXTweak.plist" "/var/jb/Library/MobileSubstrate/DynamicLibraries/ProjectXTweak.plist"
else
    JBPREFIX=""
    echo "✅ Traditional jailbreak detected, using no prefix"
    if [ -f /Library/LaunchDaemons/wangmc.projectx.plist ]; then
        launchctl load /Library/LaunchDaemons/wangmc.projectx.plist || true
    fi
    mv "/Library/ProjectXTweak.plist" "/Library/MobileSubstrate/DynamicLibraries/ProjectXTweak.plist"
fi

APP_SOURCE=""
if [ -d /Applications/ProjectX.app ]; then
    APP_SOURCE="/Applications/ProjectX.app"
elif [ -d /var/jb/Applications/ProjectX.app ]; then
    APP_SOURCE="/var/jb/Applications/ProjectX.app"
else
    APP_SOURCE="$(find /var/jb/Applications /Applications -maxdepth 3 -name ProjectX.app -print -quit 2>/dev/null)"
fi

if [ -n "$APP_SOURCE" ]; then
    if [ ! -d /Applications/ProjectX.app ]; then
        mkdir -p /Applications
        cp -R "$APP_SOURCE" /Applications/ProjectX.app
    fi
    if [ -d /var/jb/Applications ] && [ ! -d /var/jb/Applications/ProjectX.app ]; then
        cp -R "$APP_SOURCE" /var/jb/Applications/ProjectX.app
    fi
fi

if command -v uicache >/dev/null; then
    if [ -d /Applications/ProjectX.app ]; then
        uicache -p /Applications/ProjectX.app || true
    fi
    if [ -d /var/jb/Applications/ProjectX.app ]; then
        uicache -p /var/jb/Applications/ProjectX.app || true
    fi
fi

echo "✅ WeaponX installation completed!"
exit 0 
